# deploy/personal-website.service
[Unit]
Description=Personal Website (Express)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/personal-website
ExecStart=/usr/local/bin/node server/index.js
Restart=on-failure
RestartSec=5

Environment=PORT=8787
Environment=NODE_ENV=production

# 安全加固：限制进程权限
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/personal-website/public
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
